# workers/crontab
*/1 * * * * node /Users/student/2015-10-web-historian/workers/htmlfetcher.js

